# crates/l2/tee/mkosi.extra/etc/systemd/system/payload.service
[Unit]
After=network-online.target

[Service]
Type=oneshot
ExecStart=/payload
StandardOutput=journal+console
