# ansible/templates/services/ethereum-metrics-exporter.service.j2
[Unit]
Description=Ethereum Metrics Exporter
After=network.target

[Service]
ExecStart=/usr/local/bin/ethereum-metrics-exporter
Restart=on-failure

[Install]
WantedBy=multi-user.target
