# ansible/templates/services/explorer-backend.service.j2
[Unit]
Description=Explorer Backend
After=network.target

[Service]
ExecStart=

Restart=on-failure

[Install]
WantedBy=multi-user.target
